# redteam.service
[Unit]
Description=RedTeam Security Framework
After=network.target

[Service]
Type=simple
User=redteam
WorkingDirectory=/opt/redteam-framework
ExecStart=/usr/bin/python3 /opt/redteam-framework/redteam.py --config /etc/redteam/config.json
Restart=on-failure
RestartSec=60

[Install]
WantedBy=multi-user.target
